# templates/backup/mariadb-backup.timer
[Timer]
{{ mariadb.backup.timer }}

[Install]
WantedBy=timers.target
