# templates/backup/mariadb-backup.service
[Service]
{# mysqldump uses the HOME environment variable to locate .my.cnf #}
Environment=HOME=%h
ExecStart=/usr/local/bin/mariadb-backup
